# files/glances.service.j2
[Unit]
Description=Glances
After=network.target

[Service]
ExecStart=/opt/glances/bin/glances -w --password
Restart=always
RemainAfterExit=no

[Install]
WantedBy=multi-user.target